# baileys-mcp-server/scripts/whatsapp-qr-notify.timer
[Unit]
Description=Poll WhatsApp QR every minute and notify Telegram

[Timer]
OnCalendar=*:0/1  # Every minute
Persistent=true
Unit=whatsapp-qr-notify@%i.service

[Install]
WantedBy=timers.target